# conf/statsdaemon.conf
# statsite
#

description     "statsdaemon StatsD service"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30
limit nofile 65550 65550

script
  exec /opt/statsdaemon/statsdaemon -prefix="e0858dfc-0e6d-4004-8393-844a5ddda73a.stats." -graphite=carbon.hostedgraphite.com:2003 >> /mnt/log/statsdaemon.log
end script
